# containers/valve-infra/rootfs/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
#ExecStart=-/usr/bin/agetty --autologin root --noclear %I $TERM
ExecStart=-/usr/bin/su - root -c "locale-gen -a && tmuxp load dashboard"
StandardInput=tty
StandardOutput=tty
